# backend/supervisord.conf
; backend/supervisord.conf — D-17: one container == one "server", two processes.
;
; Three settings below are mandatory, not stylistic (RESEARCH Pitfall 5):
;   nodaemon=true            — else supervisord daemonises, PID 1 exits, container stops
;   nginx -g "daemon off;"   — else nginx forks and exits, supervisord restart-loops it
;   *_logfile_maxbytes=0     — else supervisord tries to rotate /dev/stdout and dies

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
